# tools/systemd/stockdash-backup-health.service
[Unit]
Description=StockDash backup health check

[Service]
Type=oneshot
User=shynwo
Environment=HOME=/home/shynwo
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/python3 /home/shynwo/stockdash/scripts/backup_monitor.py timer-check
Nice=10

