# tools/systemd/stockdash-restore-test.timer
[Unit]
Description=Run StockDash restore test every week

[Timer]
OnCalendar=Sun *-*-* 04:45:00
Persistent=true
AccuracySec=1min
Unit=stockdash-restore-test.service

[Install]
WantedBy=timers.target

